# systemd/wg-quick@wg0.service.d/10-after-modem.conf
# 让 wg-quick@wg0 严格 After 5G 拨号 + 默认路由就绪，
# 避免开机瞬间出现 "RTNETLINK answers: Network is unreachable"。
#
# 安装：
#   sudo install -Dm644 systemd/wg-quick@wg0.service.d/10-after-modem.conf \
#        /etc/systemd/system/wg-quick@wg0.service.d/10-after-modem.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart wg-quick@wg0   # 可选：让新顺序立即生效
[Unit]
After=ModemManager.service NetworkManager-wait-online.service network-online.target
Wants=ModemManager.service NetworkManager-wait-online.service network-online.target

[Service]
# 5G NSA/SA 切换可能让 endpoint DNS 解析瞬时失败：保留官方默认（infinity）即可
Environment=WG_ENDPOINT_RESOLUTION_RETRIES=infinity
# 启动前最多再等 3s 路由真正就绪，再去 add ip route（避免 unreachable 噪声）
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do ip route show default | grep -q wwan0 && exit 0; sleep 0.1; done; exit 0'
